# scripts/systemd/freeradius.service.d/sync-after-restart.conf
# Run sync-to-node2 after every start/restart of FreeRADIUS. The "+" prefix
# runs the script with full privileges (root) so it can read config and use
# /root/.ssh. Requires systemd 231+.
#
# Install:
#   sudo mkdir -p /etc/systemd/system/freeradius.service.d
#   sudo cp sync-after-restart.conf /etc/systemd/system/freeradius.service.d/
#   sudo systemctl daemon-reload
# Ensure root can SSH to Node 2: sudo ssh-copy-id user@node2

[Service]
ExecStartPost=+/usr/local/bin/sync-freeradius-to-node2.sh
